# .linux/scheduled-shutdown.timer
[Timer]
OnCalendar=20:30
Persistent=true
