# .linux/scheduled-shutdown.service
[Unit]
Description=Scheduled Shutdown at 9 PM

[Service]
Type=oneshot
ExecStart=/sbin/shutdown -h now
